# deploy/kalshiedge.service
[Unit]
Description=KalshiEdge Trading Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/KalshiEdge
EnvironmentFile=/root/KalshiEdge/.env
ExecStart=/root/KalshiEdge/.venv/bin/python -m kalshiedge.main
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
